# private_dot_config/systemd/user/xremap.service
[Unit]
Description=xremap

[Service]
Type=oneshot
KillMode=process
ExecStart=/home/riou/.cargo/bin/xremap /home/riou/.config/xremap/config.yml
ExecStop=/usr/bin/killall xremap
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
